# supervisor/breadcount-container-main.conf
[program:breadcount-container-main]
; Container tracking application (sale point / صالة)
; QR code detection and container tracking via container_main.py
command=/usr/bin/python3 container_main.py
directory=/home/sunrise/ConvuyerBreadCounting
user=sunrise
autostart=false
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stopsignal=INT
redirect_stderr=true
stdout_logfile=/home/sunrise/ConvuyerBreadCounting/data/logs/container-main.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=HOME="/home/sunrise"
