# supervisor/breadcount-container-ros2.conf
[program:breadcount-container-ros2]
; ROS2 pipeline for container camera (sale point / صالة)
; Launches hobot_rtsp_client + hobot_codec for 192.168.2.118
command=/bin/bash -c "source /opt/tros/humble/setup.bash && ros2 launch src/container/ros2/ContainerRos2Launcher.py"
directory=/home/sunrise/ConvuyerBreadCounting
user=sunrise
autostart=false
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/home/sunrise/ConvuyerBreadCounting/data/logs/container-ros2.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=HOME="/home/sunrise",RMW_IMPLEMENTATION="rmw_fastrtps_cpp",ROS_DISABLE_LOANED_MESSAGES="0"
