# supervisor/breadcount-container-main.conf
[program:breadcount-container-main]
; Container tracking application (sale point / صالة)
; QR code detection and container tracking via container_main.py
command=/bin/bash -c "source /opt/tros/humble/setup.bash && exec python3 container_main.py"
directory=/home/sunrise/ConvuyerBreadCounting
user=sunrise
autostart=true
autorestart=true
startsecs=10
startretries=5
stopwaitsecs=15
stopsignal=INT
redirect_stderr=true
stdout_logfile=/home/sunrise/ConvuyerBreadCounting/data/logs/container-main.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=HOME="/home/sunrise",RMW_IMPLEMENTATION="rmw_fastrtps_cpp",ROS_DISABLE_LOANED_MESSAGES="0"
