# supervisor/breadcount-container-main.conf
[program:breadcount-container-main]
; Container tracking application (sale point / صالة)
; QR code detection and container tracking via container_main.py
command=/bin/bash -c "source /opt/tros/humble/setup.bash && exec python3 container_main.py"
directory=/home/sunrise/ConvuyerBreadCounting
user=sunrise
autostart=true
autorestart=true
startsecs=10
startretries=5
stopwaitsecs=15
stopsignal=INT
redirect_stderr=false
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile=/tmp/container-main.stderr.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2
environment=HOME="/home/sunrise",RMW_IMPLEMENTATION="rmw_fastrtps_cpp",ROS_DISABLE_LOANED_MESSAGES="0",APP_LOG_BASENAME="container_main",APP_LOG_CONSOLE_STREAM="none"
